# Compiled/crontab
0,5,10,15,20,25,30,35,40,45,50,55 * * * * /usr/local/scripts/nationalgrid.sh #Update our database with grid stats

